# scripts/pension-notifications.service
[Unit]
Description=Pension Notification Service
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/pension
ExecStart=/usr/bin/npm run notifications
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pension-notifications
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 